# overlay/etc/init.d/headless_bootstrap
#!/sbin/openrc-run

# SPDX-FileCopyrightText: Copyright 2022-2026, macmpi
# SPDX-License-Identifier: MIT

description="Headless main boostrappring script"
name="Headless bootstrap"

command="/tmp/.ALHB/headless_bootstrap"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	use logger
    want dev-settle
    need localmount
    before net agetty.ttyGS0
    after bootmisc hwdrivers modules
}
