# overlay/etc/init.d/headless_cleanup
#!/sbin/openrc-run

# SPDX-FileCopyrightText: Copyright 2022-2025, macmpi
# SPDX-License-Identifier: MIT

description="Headless cleanup script"
name="Headless cleanup"

command="/tmp/.trash/headless_cleanup"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

