# overlay/etc/init.d/headless_bootstrap
#!/sbin/openrc-run

# SPDX-FileCopyrightText: Copyright 2022-2025, macmpi
# SPDX-License-Identifier: MIT

description="Headless main boostrappring script"
name="Headless bootstrap"

command="/tmp/.ALHB/headless_bootstrap"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

